# bowen.service
Description=bowen api
After=network.target

[Service]
User=root
Group=www-data
ExecStart=/var/www/go/bowen
[Install]
WantedBy=multi-user.target