# api/thelocr.service
Description=the locr api
After=network.target

[Service]
User=root
Group=www-data
ExecStart=/var/www/go/locr
[Install]
WantedBy=multi-user.target